# ansible/roles/lotus_fountain/templates/lotus-noncefix.timer.j2
[Unit]
Description=Lotus Fountain NONCE timer

[Timer]
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target
