# ansible/roles/lotus_fullnode/templates/lotus-bootstrap.timer.j2
[Unit]
Description=Lotus Drand Timer

[Timer]
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target
